# deploy/bayou-accounts.service
[Unit]
Description=Bayou Accounts
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bayou
Group=bayou
WorkingDirectory=/var/lib/bayou/shared
EnvironmentFile=/etc/bayou/tls.env
ExecStart=/opt/bayou/current/bin/accounts
Restart=on-failure
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/bayou/shared
ReadOnlyPaths=/etc/bayou/tls /etc/bayou/tls.env

[Install]
WantedBy=multi-user.target
