# deploy/bayou-gameserver.service
[Unit]
Description=Bayou Game Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bayou
Group=bayou
WorkingDirectory=/var/lib/bayou/shared
EnvironmentFile=/etc/bayou/tls.env
ExecStart=/opt/bayou/current/bin/gameserver --config /etc/bayou/gameserver.cfg
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/bayou/shared
ReadOnlyPaths=/etc/bayou/tls /etc/bayou/tls.env /etc/bayou/gameserver.cfg

[Install]
WantedBy=multi-user.target
