# deploy/bayou-cardserver.service
[Unit]
Description=Bayou Card Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bayou
Group=bayou
WorkingDirectory=/var/lib/bayou/shared
EnvironmentFile=/etc/bayou/tls.env
ExecStart=/opt/bayou/current/bin/cardserver
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/bayou/shared
ReadOnlyPaths=/etc/bayou/tls /etc/bayou/tls.env

[Install]
WantedBy=multi-user.target
